# root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /code

[[ $SKIP_NPM_INSTALL ]] || {
  echo "==> Install NPM dependencies"
  npm install
}

echo "==> Starting node via npm start"
exec npm start
