# root/etc/services.d/app/run
#!/bin/sh

cd /code && \
    npm install --no-bin-links && \
    npm start
